# systemd/nodemanager.service
[Unit]
Description=PlanetLab Nodemanager
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nodemanager --config /etc/edgenet/nodemanager.conf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
